# backend/deploy/shine-app-backend.service
[Unit]
Description=Shine App Backend (FastAPI)
After=network.target
# Independent of the bot — no dependency on vpn-bot.service

[Service]
User=root
WorkingDirectory=/root/shine-app/backend
EnvironmentFile=/root/shine-app/backend/.env
ExecStart=/root/shine-app/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shine-app-backend

[Install]
WantedBy=multi-user.target
